# mywebapp.socket
[Unit]
Description=Socket for MyWebApp

[Socket]
ListenStream=127.0.0.1:5000
NoDelay=true
Accept=false

[Install]
WantedBy=sockets.target